# traefik1.7-python3.7-alpine3.8/supervisord.ini
[supervisord]
nodaemon=true

[program:gunicorn]
directory=%(PWD)s
command=gunicorn %(ENV_APP_MODULE)s -c %(ENV_GUNICORN_CONF)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:traefik]
directory=%(PWD)s
command=/usr/local/bin/traefik -c %(ENV_TRAEFIK_CONFIG)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
